# salt/files/elasticsearch/postinstall.sh
/bin/systemctl daemon-reload
/bin/systemctl enable elasticsearch.service
/bin/systemctl restart elasticsearch.service
